# systemd/ttyd-bashland@.service
[Unit]
Description=ttyd bashland (%i mode)
After=network-online.target docker.service bashland-network.service
Wants=network-online.target
Requires=docker.service bashland-network.service

[Service]
Type=simple
User=ttyd
Group=ttyd
SupplementaryGroups=docker

EnvironmentFile=/etc/bashland/%i.env
Environment=HOME=/var/lib/ttyd

ExecStart=/usr/local/bin/ttyd \
  --port ${PORT} \
  --interface 127.0.0.1 \
  --writable \
  --max-clients 200 \
  --terminal-type xterm-256color \
  --base-path ${BASE_PATH} \
  -t titleFixed=BashLand \
  /opt/bashland/scripts/spawn-session.sh %i

Restart=always
RestartSec=2
LimitNOFILE=65536

ProtectSystem=strict
ReadWritePaths=/srv/bashland/logs /var/lib/ttyd
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
